# docker/api/config/cleanup.cron
30 1 * * * (cd /opt/cnaas/venv/cnaas-nac/src; source ../../bin/activate; python3 -m cnaas_nac.tools.db_cleanup; echo "Cleanup job finished") >> /var/log/cron.log
